# etc/supervisor-example.conf
[group:segagame]
programs=segagame-bot

[program:segagame-bot]
directory=/home/%USERNAME%
command=/home/%USERNAME%/venv/bin/python -m app.bot_runner
user=%USERNAME%
stopsignal=INT
autorestart=true
autostart=true
redirect_stderr=true
stdout_logfile=/home/%USERNAME%/logs/bot.txt
stdout_logfile_maxbytes=1MB
stdout_capture_maxbytes=1MB
stdout_logfile_backups=10
